# scripts/com.mygym.icloud-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mygym.icloud-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/i_kawano/Documents/crowd_data_dashboard_v2/scripts/icloud-sync.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Run every day at 00:05 -->
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>5</integer>
        </dict>
        <!-- Run every day at 12:05 -->
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>5</integer>
        </dict>
        <!-- Run every day at 18:05 -->
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>5</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/i_kawano/Documents/crowd_data_dashboard_v2/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/i_kawano/Documents/crowd_data_dashboard_v2/logs/launchd-stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>